# templates/whatsapp/launchagent_template.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zxlab.meu-agente-watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/rcastrodigital/meu-agente/watcher.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/rcastrodigital/.zxlab-mission-control/logs/meu-agente-watcher.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/rcastrodigital/.zxlab-mission-control/logs/meu-agente-watcher-error.log</string>

    <key>StartInterval</key>
    <integer>60</integer>
</dict>
</plist>
